# systemd-timers/minecraft-hourly-backup.service
[Unit]
Description="The Minecraft hourly backup"
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/minecraft/backup.sh hourly
