# systemd-timers/minecraft-monthly-backup.service
[Unit]
Description="Uhe Minecraft monthly backup"
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/minecraft/backup.sh monthly
