# systemd-timers/minecraft-daily-backup.service
[Unit]
Description="The Minecraft daily backup"
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/minecraft/backup.sh daily
